# supervisord.conf
[supervisord]
nodaemon=true

[program:avahi-daemon]
command=/usr/sbin/avahi-daemon

[program:openhab_homekit_bridge]
command=/start-ohb.sh
